# deploy/systemd/dialer.service
# systemd unit for the Visark dialer.
# Install: see deploy/README.md ("Run as a service").
#
# Adjust the three paths/user below to match your server:
#   - User/Group          : the non-root account that owns the deploy dir
#   - WorkingDirectory    : where this repo lives (must contain .env + static/)
#   - ExecStart venv path : <WorkingDirectory>/venv/bin/gunicorn
#
# NOTE: --workers 1 is deliberate. Call/message history is in-process memory,
# so a second worker would have its own copy and the UI would see inconsistent
# data. Threads give concurrency while sharing that memory. The app loads .env
# itself (python-dotenv), so no EnvironmentFile here — systemd would not strip
# the inline comments in .env and would corrupt the credentials.

[Unit]
Description=Visark cold-call dialer (WebRTC softphone)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=dialer
Group=dialer
WorkingDirectory=/opt/igdruhsil-dialer
ExecStart=/opt/igdruhsil-dialer/venv/bin/gunicorn --workers 1 --threads 4 --bind 127.0.0.1:3001 app:app

# Restart on crash AND survive reboots (enable the unit — see README)
Restart=always
RestartSec=3

# Light hardening
NoNewPrivileges=true
PrivateTmp=true

[Install]
WantedBy=multi-user.target
